# config_examples/crafty-supervisor.service
[Unit]
Description=Crafty Minecraft server supervisor
After=network.target

[Service]
Type=simple
User=crafty
WorkingDirectory=/home/crafty_fork
ExecStart=/var/opt/minecraft/crafty/.venv/bin/python -m app.classes.server_supervisor --socket /var/opt/minecraft/crafty-fork-data/supervisor.sock
Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
